# Mr-Ed/red.socket
[Socket]
ListenStream=1971
Accept=yes

[Install]
WantedBy=sockets.target
